# Send mess to telegram/check-crlCA.service
[Unit]
Description=Get CRL days to update

[Service]
Type=simple
ExecStart=/opt/scripts/check-crlCA.py

[Install]
WantedBy=multi-user.target
