# Send mess to telegram/check-letsencript.service
[Unit]
Description=Get CRL days to update

[Service]
Type=simple
ExecStart=/opt/scripts/check-letsencript.py

[Install]
WantedBy=multi-user.target
